# deploy/systemd/obs-backup.service
[Unit]
Description=Observatory daily SQLite backup
After=local-fs.target
ConditionPathIsMountPoint=/mnt/backup

[Service]
Type=oneshot
User=observatory
Group=observatory
EnvironmentFile=/etc/observatory/observatory.env
ExecStart=/opt/observatory/.venv/bin/python /opt/observatory/scripts/backup.py
StandardOutput=journal
StandardError=journal
SyslogIdentifier=obs-backup
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
