# deploy/systemd/obs-bgs-poll.service
[Unit]
Description=Observatory BGS earthquake poller
Documentation=http://earthquakes.bgs.ac.uk/
After=network-online.target chrony.service time-sync.target
Wants=network-online.target

[Service]
Type=oneshot
User=observatory
Group=observatory
EnvironmentFile=/etc/observatory/observatory.env
ExecStart=/opt/observatory/.venv/bin/python -m observatory.pollers.bgs
StandardOutput=journal
StandardError=journal
SyslogIdentifier=obs-bgs-poll
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
